# BackUPrsync/fran@192.168.18.111/run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2did-dm\x2duuid\x2dLVM\x2dyrPVKcXLl0mLBZJU2kFGCC9dcBbsHEJS7dJd7U77hRsoj9RX14Jc6C12KREe3lrn.target

[Mount]
Where=/
What=/dev/disk/by-id/dm-uuid-LVM-yrPVKcXLl0mLBZJU2kFGCC9dcBbsHEJS7dJd7U77hRsoj9RX14Jc6C12KREe3lrn
Type=ext4
